# update-rkn.timer
[Unit]
Description=Update RKN list hourly
 
[Timer]
OnCalendar=hourly
Persistent=true
 
[Install]
WantedBy=timers.target
